# uninstall/etc/udev/rules.d/uninstall_ksm09.rules
KERNEL=="loop0", RUN+="/root/uninstall_ksm09.sh"